# provision/ansible/roles/vllm_colpali/templates/vllm-colpali.service.j2
[Unit]
Description=vLLM ColPali Vision Model Server (GPU {{ vllm_colpali_cuda_visible_devices }})
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ vllm_colpali_user }}
Group={{ vllm_colpali_group }}
WorkingDirectory={{ vllm_colpali_home }}

# Environment
Environment="CUDA_VISIBLE_DEVICES={{ vllm_colpali_cuda_visible_devices }}"
Environment="HF_HOME={{ vllm_colpali_cache_dir }}"
Environment="VLLM_CACHE_ROOT={{ vllm_colpali_home }}/.cache/vllm"
Environment="TRANSFORMERS_CACHE={{ vllm_colpali_cache_dir }}"
Environment="HF_HUB_CACHE={{ vllm_colpali_cache_dir }}/hub"
{% if vllm_colpali_hf_token %}
Environment="HF_TOKEN={{ vllm_colpali_hf_token }}"
{% endif %}

# vLLM Command with LoRA support for ColPali
# ColPali is a LoRA adapter on top of PaliGemma base model
ExecStart={{ vllm_colpali_venv_dir }}/bin/python -m vllm.entrypoints.openai.api_server \
    --model google/paligemma-3b-pt-448 \
    --host {{ vllm_colpali_host }} \
    --port {{ vllm_colpali_port }} \
    --tensor-parallel-size {{ vllm_colpali_tensor_parallel_size }} \
    --gpu-memory-utilization {{ vllm_colpali_gpu_memory_utilization }} \
    --max-model-len {{ vllm_colpali_max_model_len }} \
    --max-num-seqs {{ vllm_colpali_max_num_seqs }} \
    --served-model-name {{ vllm_colpali_served_model_name }} \
    --trust-remote-code \
    --enable-lora \
    --max-loras 1 \
    --max-lora-rank 32 \
    --lora-modules colpali={{ vllm_colpali_model }} \
    --limit-mm-per-prompt '{"image": 1}' \
{% if not vllm_colpali_enable_prefix_caching %}
    --disable-frontend-multiprocessing \
{% endif %}
{% if vllm_colpali_disable_log_stats %}
    --disable-log-stats \
{% endif %}
    --dtype auto

# Resource Limits
MemoryLimit={{ vllm_colpali_memory_limit }}
CPUQuota={{ vllm_colpali_cpu_limit | int * 100 }}%

# Restart Configuration
Restart={% if vllm_colpali_restart_on_failure %}on-failure{% else %}no{% endif %}
RestartSec={{ vllm_colpali_restart_sec }}
TimeoutStartSec={{ vllm_colpali_startup_timeout }}
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm-colpali

[Install]
WantedBy=multi-user.target

