# provision/docker/supervisord.conf
# =============================================================================
# Supervisord Configuration for Core Apps Container
# =============================================================================
#
# Manages three processes:
#   1. busibox-portal - Next.js application on port 3000
#   2. busibox-agents - Next.js application on port 3001
#   3. busibox-appbuilder - Next.js application on port 3004
#
# Process management:
#   - Apps start as autostart=false, started by entrypoint after deployment
#   - Use supervisorctl to manage individual processes
#
# Commands:
#   supervisorctl status           - Show all process status
#   supervisorctl restart busibox-portal - Restart Busibox Portal
#   supervisorctl stop busibox-agents - Stop Agent Manager
#   supervisorctl tail -f busibox-portal - Follow Busibox Portal logs
#
# =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root
loglevel=info

# =============================================================================
# Busibox Portal - Next.js Application
# =============================================================================
# Starts with autostart=false - entrypoint script starts after deployment
# Uses standalone output mode for minimal memory footprint (~300-500MB vs ~3.5GB)
# autorestart=true - restart automatically if the process crashes
[program:busibox-portal]
command=node apps/portal/server.js
directory=/srv/busibox-portal/.next/standalone
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/portal",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Agent Manager - Next.js Application
# =============================================================================
# Starts with autostart=false - entrypoint script starts after deployment
# Uses standalone output mode for minimal memory footprint (~200-400MB vs ~3.5GB)
# autorestart=true - restart automatically if the process crashes
[program:busibox-agents]
command=node apps/agents/server.js
directory=/srv/busibox-agents/.next/standalone
environment=PORT="3001",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/agents",DEFAULT_API_AUDIENCE="agent-api",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Admin - Next.js Application
# =============================================================================
[program:busibox-admin]
command=node apps/admin/server.js
directory=/srv/busibox-admin/.next/standalone
environment=PORT="3002",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/admin",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Chat - Next.js Application
# =============================================================================
[program:busibox-chat]
command=node apps/chat/server.js
directory=/srv/busibox-chat/.next/standalone
environment=PORT="3003",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/chat",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# App Builder - Next.js Application
# =============================================================================
# Starts with autostart=false - entrypoint script starts after deployment
# Uses standalone output mode for minimal memory footprint
# autorestart=true - restart automatically if the process crashes
[program:busibox-appbuilder]
command=node apps/appbuilder/server.js
directory=/srv/busibox-appbuilder/.next/standalone
environment=PORT="3004",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/builder",APP_NAME="busibox-appbuilder",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Media - Next.js Application
# =============================================================================
[program:busibox-media]
command=node apps/media/server.js
directory=/srv/busibox-media/.next/standalone
environment=PORT="3005",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/media",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Documents - Next.js Application
# =============================================================================
[program:busibox-documents]
command=node apps/documents/server.js
directory=/srv/busibox-documents/.next/standalone
environment=PORT="3006",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/documents",NODE_OPTIONS="--max-old-space-size=512"
autostart=false
autorestart=true
priority=20
startsecs=10
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Unix Socket for supervisorctl
# =============================================================================
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# =============================================================================
# Event Listener - Forward stdout/stderr to Docker logs
# =============================================================================
[eventlistener:stdout]
command=cat
buffer_size=100
events=PROCESS_LOG
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
