# provision/docker/supervisord.conf
# =============================================================================
# Supervisord Configuration for Core Apps Container
# =============================================================================
#
# Manages three processes:
#   1. nginx - Reverse proxy (always running)
#   2. ai-portal - Next.js application on port 3000
#   3. agent-manager - Next.js application on port 3001
#
# Process management:
#   - nginx starts immediately (autostart=true)
#   - Apps start as autostart=false, started by entrypoint after deployment
#   - Use supervisorctl to manage individual processes
#
# Commands:
#   supervisorctl status           - Show all process status
#   supervisorctl restart nginx    - Restart nginx
#   supervisorctl restart ai-portal - Restart AI Portal
#   supervisorctl stop agent-manager - Stop Agent Manager
#   supervisorctl tail -f ai-portal - Follow AI Portal logs
#
# =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root
loglevel=info

# =============================================================================
# Nginx - Reverse Proxy
# =============================================================================
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# AI Portal - Next.js Application
# =============================================================================
# Starts with autostart=false - entrypoint script starts after deployment
# Uses npm start (next start) since standalone mode is not enabled
# autorestart=false - don't restart on failure, just stop
[program:ai-portal]
command=npm start
directory=/srv/ai-portal
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/portal"
autostart=false
autorestart=false
priority=20
startsecs=10
startretries=1
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Agent Manager - Next.js Application
# =============================================================================
# Starts with autostart=false - entrypoint script starts after deployment
# Uses npm start (next start) since standalone mode is not enabled
# autorestart=false - don't restart on failure, just stop
[program:agent-manager]
command=npm start
directory=/srv/agent-manager
environment=PORT="3001",HOSTNAME="0.0.0.0",NODE_ENV="production",NEXT_PUBLIC_BASE_PATH="/agents"
autostart=false
autorestart=false
priority=20
startsecs=10
startretries=1
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# =============================================================================
# Unix Socket for supervisorctl
# =============================================================================
[unix_http_server]
file=/run/supervisord.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# =============================================================================
# Event Listener - Forward stdout/stderr to Docker logs
# =============================================================================
[eventlistener:stdout]
command=cat
buffer_size=100
events=PROCESS_LOG
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
