# provision/ansible/roles/ingest/templates/ingest-api.service.j2
[Unit]
Description=Busibox Ingestion API Server
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/srv/ingest/src
EnvironmentFile=/srv/ingest/.env
ExecStart=/srv/ingest/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ingest-api

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

