# provision/ansible/roles/deploy_api/templates/deploy-api.service.j2
[Unit]
Description=Busibox Deployment Service
After=network.target authz-api.service
Wants=authz-api.service

[Service]
User=root
Group=root
WorkingDirectory={{ deploy_api_app_dir }}
EnvironmentFile={{ deploy_api_base_dir }}/.env
# Run from app dir so 'src' is a proper package with relative imports
ExecStart={{ deploy_api_venv_dir }}/bin/uvicorn src.main:app --host 0.0.0.0 --port {{ deploy_api_service_port }}
Restart=always
RestartSec=5

# Security hardening (but needs root for SSH/Ansible)
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
