# provision/ansible/roles/vllm_embedding/templates/vllm-embedding.service.j2
[Unit]
Description=vLLM Embedding Server - {{ vllm_embedding_model }}
Documentation=https://docs.vllm.ai/
After=network.target

[Service]
Type=simple
User={{ vllm_embedding_user }}
Group={{ vllm_embedding_group }}
WorkingDirectory={{ vllm_embedding_home }}

# Environment configuration
Environment="CUDA_VISIBLE_DEVICES={{ vllm_embedding_cuda_visible_devices }}"
Environment="HF_HOME={{ vllm_embedding_cache_dir }}"
{% if vllm_embedding_hf_token %}
Environment="HF_TOKEN={{ vllm_embedding_hf_token }}"
{% endif %}

# Service execution
ExecStart={{ vllm_embedding_venv_dir }}/bin/python -m vllm.entrypoints.openai.api_server \
  --model {{ vllm_embedding_model }} \
  --host {{ vllm_embedding_host }} \
  --port {{ vllm_embedding_port }} \
  --tensor-parallel-size {{ vllm_embedding_tensor_parallel_size }} \
  --gpu-memory-utilization {{ vllm_embedding_gpu_memory_utilization }} \
  --max-model-len {{ vllm_embedding_max_model_len }} \
  --max-num-seqs {{ vllm_embedding_max_num_seqs }} \
{% if vllm_embedding_max_num_batched_tokens %}
  --max-num-batched-tokens {{ vllm_embedding_max_num_batched_tokens }} \
{% endif %}
{% if vllm_embedding_enable_prefix_caching %}
  --enable-prefix-caching \
{% endif %}
{% if vllm_embedding_trust_remote_code %}
  --trust-remote-code \
{% endif %}
{% if vllm_embedding_served_model_name %}
  --served-model-name {{ vllm_embedding_served_model_name }} \
{% endif %}
{% if vllm_embedding_enable_api_key and vllm_embedding_api_key %}
  --api-key {{ vllm_embedding_api_key }} \
{% endif %}
{% if vllm_embedding_disable_log_stats %}
  --disable-log-stats
{% endif %}

# Restart configuration
Restart={{ 'always' if vllm_embedding_restart_on_failure else 'no' }}
RestartSec={{ vllm_embedding_restart_sec }}
TimeoutStartSec={{ vllm_embedding_startup_timeout }}

# Resource limits
{% if vllm_embedding_memory_limit is defined %}
MemoryMax={{ vllm_embedding_memory_limit }}
{% endif %}
{% if vllm_embedding_cpu_limit is defined %}
CPUQuota={{ vllm_embedding_cpu_limit }}00%
{% endif %}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vllm_embedding_home }}
ReadWritePaths={{ vllm_embedding_cache_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm-embedding

[Install]
WantedBy=multi-user.target

