# provision/ansible/roles/vllm_8005/templates/vllm.service.j2
[Unit]
Description=vLLM Model Server (Port {{ vllm_port }}) - {{ vllm_default_model }}
Documentation=https://docs.vllm.ai/
After=network.target

[Service]
Type=simple
User={{ vllm_user }}
Group={{ vllm_group }}
WorkingDirectory={{ vllm_home }}

# Environment configuration
EnvironmentFile=/etc/default/{{ vllm_service_name }}

# Service execution
ExecStart={{ vllm_venv_dir }}/bin/python -m vllm.entrypoints.openai.api_server \
{% if vllm_enable_lora and vllm_lora_base_model %}
  --model {{ vllm_lora_base_model }} \
  --enable-lora \
  --lora-modules colpali={{ vllm_default_model }} \
  --max-loras {{ vllm_max_loras }} \
  --max-lora-rank {{ vllm_max_lora_rank }} \
{% else %}
  --model {{ vllm_default_model }} \
{% endif %}
  --host {{ vllm_host }} \
  --port {{ vllm_port }} \
  --tensor-parallel-size {{ vllm_tensor_parallel_size }} \
  --gpu-memory-utilization {{ vllm_gpu_memory_utilization }} \
  --max-model-len {{ vllm_max_model_len }} \
  --max-num-seqs {{ vllm_max_num_seqs }} \
{% if vllm_cpu_offload_gb is defined and vllm_cpu_offload_gb > 0 %}
  --cpu-offload-gb {{ vllm_cpu_offload_gb }} \
{% endif %}
{% if vllm_max_num_batched_tokens %}
  --max-num-batched-tokens {{ vllm_max_num_batched_tokens }} \
{% endif %}
{% if vllm_enable_prefix_caching %}
  --enable-prefix-caching \
{% endif %}
{% if vllm_enable_chunked_prefill is defined and vllm_enable_chunked_prefill %}
  --enable-chunked-prefill \
{% endif %}
{% if vllm_enable_expert_parallel is defined and vllm_enable_expert_parallel %}
  --enable-expert-parallel \
{% endif %}
{% if vllm_kv_cache_dtype is defined and vllm_kv_cache_dtype | length > 0 %}
  --kv-cache-dtype {{ vllm_kv_cache_dtype }} \
{% endif %}
{% if vllm_reasoning_parser is defined and vllm_reasoning_parser | length > 0 %}
  --enable-reasoning --reasoning-parser {{ vllm_reasoning_parser }} \
{% endif %}
{% if vllm_trust_remote_code %}
  --trust-remote-code \
{% endif %}
{% if vllm_served_model_name %}
  --served-model-name {{ vllm_served_model_name }} \
{% endif %}
{% if vllm_enable_api_key and vllm_api_key %}
  --api-key {{ vllm_api_key }} \
{% endif %}
{% if vllm_enable_tool_calling is defined and vllm_enable_tool_calling %}
  --enable-auto-tool-choice \
{% if vllm_tool_call_parser is defined and vllm_tool_call_parser | length > 0 %}
  --tool-call-parser {{ vllm_tool_call_parser }} \
{% endif %}
{% if vllm_tool_chat_template is defined and vllm_tool_chat_template | length > 0 %}
  --chat-template {{ vllm_tool_chat_template }} \
{% endif %}
{% endif %}
{% if vllm_guided_decoding_backend is defined and vllm_guided_decoding_backend | length > 0 %}
  --guided-decoding-backend {{ vllm_guided_decoding_backend }} \
{% endif %}
{% if vllm_multimodal is defined and vllm_multimodal %}
{% if vllm_limit_mm_per_prompt is defined and vllm_limit_mm_per_prompt | length > 0 %}
  --limit-mm-per-prompt '{{ vllm_limit_mm_per_prompt }}' \
{% endif %}
{% endif %}
{% if vllm_hf_overrides is defined and vllm_hf_overrides | length > 0 %}
  --hf-overrides '{{ vllm_hf_overrides }}' \
{% endif %}
{% if vllm_disable_log_stats %}
  --disable-log-stats
{% endif %}

# Restart configuration
Restart={{ 'always' if vllm_restart_on_failure else 'no' }}
RestartSec={{ vllm_restart_sec }}
TimeoutStartSec={{ vllm_startup_timeout }}

# Resource limits
{% if vllm_memory_limit is defined %}
MemoryMax={{ vllm_memory_limit }}
{% endif %}
{% if vllm_cpu_limit is defined %}
CPUQuota={{ vllm_cpu_limit }}00%
{% endif %}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vllm_home }}
ReadWritePaths={{ vllm_cache_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ vllm_service_name }}

[Install]
WantedBy=multi-user.target

