# provision/ansible/roles/signal_bot/templates/signal-bot.service.j2
# Signal Bot systemd Service
# Generated by Ansible - Do not edit manually
# {{ ansible_managed }}

[Unit]
Description=Signal Bot for AI Portal
Documentation=https://github.com/your-org/busibox
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ signal_bot_user }}
Group={{ signal_bot_group }}
WorkingDirectory={{ signal_bot_dir }}
EnvironmentFile={{ signal_bot_dir }}/.env

ExecStart={{ signal_bot_venv }}/bin/python -m app.main

Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ signal_bot_service_name }}

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ signal_bot_dir }}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
