# provision/ansible/roles/app_deployer/templates/app-systemd.service.j2
# {{ ansible_managed }}
# Systemd service for {{ app.name }}
# Generated by app_deployer role

[Unit]
Description={{ app.name }} Application
After=network.target
{% if app.requires is defined %}
{% for req in app.requires %}
Requires={{ req }}.service
After={{ req }}.service
{% endfor %}
{% endif %}

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{ app.deploy_path }}

# Environment variables
{% if app.env is defined %}
{% for key, value in app.env.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}

# Port configuration
Environment="PORT={{ app.port }}"
Environment="NODE_ENV=production"

# Start command
{% if app.start_command is defined %}
ExecStart=/bin/bash -c '{{ app.start_command }}'
{% else %}
ExecStart=/usr/bin/npm start -- -p {{ app.port }}
{% endif %}

# Restart configuration
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Logging - use journal for consistent timestamps
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app.name }}

# Security hardening (optional, can be enabled per-app)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target

